# deploy/systemd/trader-scrape-ebay.timer
[Unit]
Description=Run trader eBay scrape daily

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
