# deploy/systemd/trader-scrape-subito.service
[Unit]
Description=Trader scrape Subito
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=trader
Group=trader
ExecStart=/opt/trader/app/deploy/server_job.sh scrape-subito
