# deploy/systemd/trader-ai-cascade.service
[Unit]
Description=Trader AI cascade classify pending review ads
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=trader
Group=trader
EnvironmentFile=/etc/trader/trader.env
WorkingDirectory=/opt/trader/app
ExecStart=/opt/trader/app/deploy/server_job.sh ai-cascade
