# deploy/systemd/trader-backup.service
[Unit]
Description=Trader SQLite backup

[Service]
Type=oneshot
User=trader
Group=trader
ExecStart=/opt/trader/app/deploy/backup_db.sh
